# roles/homebrew/templates/brew-update.timer.j2
# {{ ansible_managed }}
[Unit]
Description=Timer for updating Brew binary
Wants=network-online.target

[Timer]
OnBootSec=10min
OnUnitInactiveSec=8h
Persistent=true

[Install]
WantedBy=timers.target
